# tests/fixtures/macos/System/Library/LaunchDaemons/com.example.beta.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.example.beta</string>
	<key>MachServices</key>
	<dict>
		<key>com.example.beta</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/betad</string>
		<string>--now</string>
	</array>
	<key>SandboxProfile</key>
	<string>com.example.beta</string>
	<key>UserName</key>
	<string>_beta</string>
</dict>
</plist>
