# tests/fixtures/linux/lib/systemd/system/test.service
[Unit]
Description=Test
[Service]
ExecStart=/usr/bin/test_daemon
